# mitra-monitor.service
[Unit]
Description=Mitra Application Health Monitor and Auto-Restart
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/bin/bash /mnt/mitradisk/app/hp/monitor_mitra.sh
Restart=always
RestartSec=10
User=root
# Environment variables can be added here if necessary
# Environment=VAR=value

[Install]
WantedBy=multi-user.target